# 50-usb_custom.rules
KERNEL=="sd*1", ACTION=="add", SUBSYSTEMS=="block", RUN+="/home/pi/Scripts/getitnow.sh %k"
